# supervisor-svgparser.conf
[program:app-uwsgi]
command = /home/svg-parser/env/bin/uwsgi --ini /home/svg-parser/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
